# supervisor/conf.d/beanstalkd.conf
[program:beanstalkd]
command = /usr/bin/beanstalkd -l 10.104.95.66 -c
stdout_logfile=/data/logs/supervisor-beanstalkd-stdout.log
stderr_logfile=/data/logs/supervisor-beanstalkd-stderr.log
